# openrc/service
#!/sbin/openrc-run

description="Updates Google Domains Dynamic DNS with current IP Address"

command="${RC_SVCNAME}"
# command_background=true

start() {
  env
  # env start-stop-daemon "${RC_SVCNAME}" \
  # --make-pidfile \
  # --pidfile /run/${RC_SVCNAME}.pid \
  # --start \
  # --background \
  # -v
}
