# openrc/service
#!/sbin/openrc-run

description="Updates Google Domains Dynamic DNS with current IP Address"

# command=/usr/local/bin/google_dynamic_dns_update_service
command=env
command_args=google_dynamic_dns_update_service
# command_background=true
